# unified-api.service
[Unit]
Description=Unified WhatsApp and CEP Validation API
After=network.target

[Service]
User=brpl
Group=brpl
WorkingDirectory=/home/brpl/wap-is-whats-number
ExecStart=/bin/bash start.sh
Restart=always
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=unified-api
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
